# cli/systemd/ethan.service
[Unit]
Description=ETHAN Cognitive OS — backend service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/fatsio/AI/Ethan

ExecStart=/usr/bin/env bash -lc 'if command -v docker >/dev/null 2>&1; then docker compose -f docker-compose.yml up -d || true; fi'
ExecStopPost=/usr/bin/env bash -lc 'if command -v docker >/dev/null 2>&1; then docker compose -f docker-compose.yml down || true; fi'

Restart=on-failure
RestartSec=10

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=default.target
