# infrastructure/systemd/ethan-core.service
[Unit]
Description=ETHAN Core AI Engine
Documentation=https://github.com/seclib/Ethan/docs
Wants=network-online.target nats.service redis.service postgres.service
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ethan-core
Restart=on-failure
RestartSec=5s
User=ethan
Group=ethan

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/ethan /var/log/ethan /var/cache/ethan

# Resources
MemoryMax=2G
CPUQuota=200%
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ethan-core

[Install]
WantedBy=multi-user.target